# files/etc/network/if-up.d/oz
#!/bin/bash
set -e

if [[ "$METHOD" = loopback ]] || [[ "$LOGICAL" = "oz0" ]] || [[ "$LOGICAL" == veth* ]] || [[ "$LOGICAL" == voz* ]]; then
	exit 0
fi

killall -SIGUSR2 oz-daemon || exit 0
